# systemd/rclone-hybrid-web.service
[Unit]
Description=Rclone Hybrid Web API (queue + scheduler + watcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rclone-hybrid/hybrid/backend
EnvironmentFile=-/opt/rclone-hybrid/hybrid/.env
ExecStart=/opt/rclone-hybrid/hybrid/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
